# com.hermes.voice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.voice</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/HermesVoice.app/Contents/MacOS/HermesVoice</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/fahmidhasantaohid/Library/Logs/HermesVoice/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fahmidhasantaohid/Library/Logs/HermesVoice/stderr.log</string>
</dict>
</plist>
